# access-point/wifi-connection-monitor.service
[Unit]
Description=WiFi Connection Monitor - Fallback to AP mode if connection fails
After=network.target wpa_supplicant.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wifi-connection-monitor.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
